# profile.d/jupyter.sh
VIRTUAL_ENV=/jupyter
PATH="$VIRTUAL_ENV/bin:$PATH"
PS1="[\A] \H:\w $ "
export PS1
LD_LIBRARY_PATH=/usr/local/cuda/extras/CUPTI/lib64:/usr/local/cuda/lib64
export LD_LIBRARY_PATH
PATH="$PATH:/usr/local/cuda/bin:/usr/local/bin:/usr/bin:/bin:/sbin:/usr/sbin:/usr/local/sbin"
export PATH
